# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = cron -f

[program:php]
command = docker-php-entrypoint php-fpm
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0